# mkosi/machines/booter/entrypointd-handler.service
[Unit]
Description=Start entrypointd portable

[Service]
Type=oneshot
ExecStart=portablectl attach --enable --now entrypointd
