# mkosi/machines/booter/mkosi.extra/usr/local/lib/systemd/system/entrypoint-watch.path
[Unit]
Description=Await entrypointd portable service

[Path]
PathExists=/var/lib/portables/entrypointd
Unit=entrypointd-handler.service

[Install]
WantedBy=multi-user.target
